# templates/update-motd.d/99-custom
#!/usr/bin/env bash
# Custom MOTD (Korean) - server status summary
set -Eeuo pipefail

# --- helpers ---------------------------------------------------------------

get_login_user() {
  # sudo/suì—ë„ ê°•í•œ "ì›ëž˜ ë¡œê·¸ì¸ ì‚¬ìš©ìž"
  logname 2>/dev/null || echo "unknown"
}

get_client_ip() {
  # PAM/update-motd í™˜ê²½ì—ì„œë„ ê°€ìž¥ ìž˜ ìž¡íž˜
  local ip
  ip="$(who -u am i 2>/dev/null | awk '{print $NF}' | tr -d '()' || true)"
  if [[ -n "${ip:-}" && "$ip" != "-" ]]; then
    echo "$ip"
    return
  fi

  # fallback: SSH í™˜ê²½ë³€ìˆ˜
  if [[ -n "${SSH_CLIENT:-}" ]]; then
    echo "${SSH_CLIENT%% *}"
    return
  fi
  if [[ -n "${SSH_CONNECTION:-}" ]]; then
    echo "${SSH_CONNECTION%% *}"
    return
  fi

  echo "(ì½˜ì†”/ë¡œì»¬)"
}

get_mem_available_human() {
  # /proc/meminfo MemAvailable ê¸°ë°˜
  awk '
    /MemAvailable:/ {
      kb=$2
      mib=kb/1024
      gib=mib/1024
      if (gib>=1) printf "%.1fGiB", gib
      else printf "%.0fMiB", mib
      exit
    }
  ' /proc/meminfo 2>/dev/null || echo "N/A"
}

get_disk_free_root() {
  df -h / 2>/dev/null | awk 'NR==2 {print $4}' || echo "N/A"
}

# --- main ------------------------------------------------------------------

LOGIN_USER="$(get_login_user)"
RUN_AS="$(whoami)"
CLIENT_IP="$(get_client_ip)"
DISK_FREE="$(get_disk_free_root)"
MEM_AVAIL="$(get_mem_available_human)"

echo "========================================"
echo " ðŸ–¥  ì„œë²„ ìƒíƒœ ìš”ì•½"
echo "----------------------------------------"
echo " ðŸ“… í˜„ìž¬ ì‹œê°„  : $(date '+%Y-%m-%d %H:%M:%S')"

if [[ "$LOGIN_USER" == "$RUN_AS" ]]; then
  echo " ðŸ‘¤ ë¡œê·¸ì¸ ê³„ì •: ${LOGIN_USER}"
else
  echo " ðŸ‘¤ ë¡œê·¸ì¸ ê³„ì •: ${LOGIN_USER} (í˜„ìž¬ ê¶Œí•œ: ${RUN_AS})"
fi

echo " ðŸŒ ì ‘ì† IP    : ${CLIENT_IP}"
echo " â±  ì—…íƒ€ìž„     : $(uptime -p 2>/dev/null || echo 'N/A')"
echo " ðŸ’¾ ë””ìŠ¤í¬(/)  : ${DISK_FREE} ì—¬ìœ "
echo " ðŸ§  ë©”ëª¨ë¦¬     : ${MEM_AVAIL} ì—¬ìœ "

# Fail2Ban (root ê¶Œí•œì¼ ë•Œë§Œ)
if [[ "$(id -u)" -eq 0 ]]; then
  if command -v fail2ban-client >/dev/null 2>&1; then
    BANNED="$(timeout 1 fail2ban-client status sshd 2>/dev/null | awk -F': ' '/Currently banned/ {print $2}')"
    [[ -n "${BANNED:-}" ]] || BANNED="0"
    echo " ðŸ” Fail2Ban   : ${BANNED} IP ì°¨ë‹¨ ì¤‘"

    if [[ "$BANNED" != "0" ]]; then
      IPS="$(timeout 1 fail2ban-client get sshd banip 2>/dev/null | tr -s ' ' | sed 's/^ //;s/ $//')"
      if [[ -n "${IPS:-}" ]]; then
        COUNT="$(echo "$IPS" | wc -w | tr -d ' ')"
        if [[ "$COUNT" -le 8 ]]; then
          echo " ðŸ§± ì°¨ë‹¨ IP    : ${IPS}"
        else
          HEAD8="$(echo "$IPS" | awk '{for(i=1;i<=8;i++) printf "%s%s",$i,(i<8?" ":"");}')"
          echo " ðŸ§± ì°¨ë‹¨ IP    : ${HEAD8} ... (+$((COUNT-8)) more)"
        fi
      fi
    fi
  else
    echo " ðŸ” Fail2Ban   : (ë¯¸ì„¤ì¹˜)"
  fi
else
  echo " ðŸ” Fail2Ban   : (ê¶Œí•œ ì—†ìŒ)"
fi

echo "========================================"
