# templates/update-motd.d/99-custom
#!/bin/bash
# Custom MOTD (Korean) - server status summary

get_client_ip() {
  if [ -n "${SSH_CLIENT:-}" ]; then
    echo "${SSH_CLIENT%% *}"
    return
  fi
  if [ -n "${SSH_CONNECTION:-}" ]; then
    echo "${SSH_CONNECTION%% *}"
    return
  fi
  if [ -n "${LOGIN_IP:-}" ]; then
    echo "${LOGIN_IP}"
    return
  fi
  echo "(ì½˜ì†”/ë¡œì»¬)"
}

get_login_user() {
  # SSH ë¡œê·¸ì¸ ì‚¬ìš©ìž: sudo/suë¡œ rootê°€ ë˜ì–´ë„ ì›ëž˜ ì‚¬ìš©ìž ìœ ì§€
  if [ -n "${SUDO_USER:-}" ]; then
    echo "${SUDO_USER}"
    return
  fi
  if [ -n "${LOGNAME:-}" ]; then
    echo "${LOGNAME}"
    return
  fi
  if [ -n "${USER:-}" ]; then
    echo "${USER}"
    return
  fi
  whoami
}

# SSH ì ‘ì† IP ìºì‹œ (sudo/suë¡œ ë“¤ì–´ì˜¬ ë•Œë¥¼ ëŒ€ë¹„)
if [ -n "${SSH_CLIENT:-}" ]; then
  export LOGIN_IP="${SSH_CLIENT%% *}"
fi

CLIENT_IP="$(get_client_ip)"

LOGIN_USER="$(get_login_user)"
RUN_AS="$(whoami)"

echo "========================================"
echo " ðŸ–¥  ì„œë²„ ìƒíƒœ ìš”ì•½"
echo "----------------------------------------"
echo " ðŸ“… í˜„ìž¬ ì‹œê°„  : $(date '+%Y-%m-%d %H:%M:%S')"

# ë¡œê·¸ì¸ ì‚¬ìš©ìžì™€ ì‹¤í–‰ ì‚¬ìš©ìžê°€ ë‹¤ë¥´ë©´ ë‘˜ ë‹¤ ë³´ì—¬ì£¼ê¸°
if [ "$LOGIN_USER" = "$RUN_AS" ]; then
  echo " ðŸ‘¤ ë¡œê·¸ì¸ ê³„ì •: ${LOGIN_USER}"
else
  echo " ðŸ‘¤ ë¡œê·¸ì¸ ê³„ì •: ${LOGIN_USER} (í˜„ìž¬ ê¶Œí•œ: ${RUN_AS})"
fi

echo " ðŸŒ ì ‘ì† IP    : ${CLIENT_IP}"
echo " â±  ì—…íƒ€ìž„     : $(uptime -p)"
echo " ðŸ’¾ ë””ìŠ¤í¬(/)  : $(df -h / | awk 'NR==2 {print $4}') ì—¬ìœ "
echo " ðŸ§  ë©”ëª¨ë¦¬     : $(free -h | awk '/^Mem:/ {print $7}') ì—¬ìœ "

# Fail2Ban ìƒíƒœ + ì°¨ë‹¨ IP ëª©ë¡
if [ "$(id -u)" -eq 0 ]; then
  BANNED="$(timeout 1 fail2ban-client status sshd 2>/dev/null | awk -F': ' '/Currently banned/ {print $2}')"
  [ -n "$BANNED" ] || BANNED="0"
  echo " ðŸ” Fail2Ban   : ${BANNED} IP ì°¨ë‹¨ ì¤‘"

  if [ "$BANNED" != "0" ]; then
    # ì°¨ë‹¨ IP ëª©ë¡(ê°€ëŠ¥í•˜ë©´)
    IPS="$(timeout 1 fail2ban-client get sshd banip 2>/dev/null | tr -s ' ' | sed 's/^ //;s/ $//')"
    if [ -n "$IPS" ]; then
      # ë„ˆë¬´ ê¸¸ë©´ ì•žë¶€ë¶„ë§Œ í‘œì‹œ (ìµœëŒ€ 6ê°œ)
      COUNT="$(echo "$IPS" | wc -w | tr -d ' ')"
      if [ "$COUNT" -le 6 ]; then
        echo " ðŸ§± ì°¨ë‹¨ IP    : ${IPS}"
      else
        HEAD6="$(echo "$IPS" | awk '{for(i=1;i<=6;i++) printf "%s%s",$i,(i<6?" ":"");}')"
        echo " ðŸ§± ì°¨ë‹¨ IP    : ${HEAD6} ... (+$((COUNT-6)) more)"
      fi
    else
      echo " ðŸ§± ì°¨ë‹¨ IP    : (ëª©ë¡ ì¡°íšŒ ë¶ˆê°€)"
    fi
  fi
else
  echo " ðŸ” Fail2Ban   : (ê¶Œí•œ ì—†ìŒ)"
fi

echo "========================================"
